# systemd/borgmatic-wavelog.timer
[Unit]
Description=Run Borgmatic Backup for wavelog

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true

[Install]
WantedBy=timers.target
