# systemd/wavelog-ping.timer
[Unit]
Description=Check if Wavelog is online

[Timer]
OnCalendar=00/1:05
Persistent=true

[Install]
WantedBy=timers.target
